# symlink/~/.config/systemd/user/niri-session.target
[Unit]
Description=Niri Session Target
Requires=graphical-session.target
After=graphical-session.target

[Install]
WantedBy=default.target
